# rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
nginx -g "pid /tmp/nginx.pid; daemon off;"